# deploy/ansible/roles/agent/files/systemd/system/ip-route-add@agent-remote.service.d/override.conf
[Unit]
After=ip-route-add@vrf10.service
After=ip-route-add@end0.service

[Service]
Environment="DESTINATION_SUBNET=10.11.0.0/16"
