# deploy/ansible/roles/agent/files/systemd/system/ip-route-add@vrf10.service.d/override.conf
[Service]
# VRF -> LINK
Environment=VRF=%i
Environment="GATEWAY=end0"
Environment="DESTINATION_SUBNET=10.50.0.0/24"
