# deploy/ansible/roles/agent/files/systemd/system/ip-route-add@end0.service.d/override.conf
[Service]
# VRF <- LINK
Environment="DESTINATION_SUBNET=10.16.0.0/16"
ExecStart=
ExecStart=/usr/local/sbin/%p.sh ${DESTINATION_SUBNET} ${VRF}
